# vconf.manifest
<manifest>
	<request>
		<domain name="_"/>
	</request>
	<assign>
		<filesystem path="/bin/*" exec_label="none" />
		<filesystem path="/usr/bin/*" exec_label="none" />
		<filesystem path="/etc/rc.d/init.d/vconf-init" exec_label="none" />
	</assign>
</manifest>
